# lib/logicblox/lib/systemd/user/lb.target
[Unit]
Description=All LogicBlox Services
Wants=lb-server.socket lb-server-admin.socket lb-compiler.service
